# docker/supervisord.conf
; ═══════════════════════════════════════════════════════════════════════════════
; Supervisor Configuration for AMTTP Platform
; Manages all services in a single container
; ═══════════════════════════════════════════════════════════════════════════════

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
childlogdir=/var/log/supervisor

; ═══════════════════════════════════════════════════════════════════════════════
; NGINX - Main Web Server & Reverse Proxy
; ═══════════════════════════════════════════════════════════════════════════════

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log

; ═══════════════════════════════════════════════════════════════════════════════
; Next.js Dashboard (Standalone Mode)
; ═══════════════════════════════════════════════════════════════════════════════

[program:nextjs]
command=node server.js
directory=/app/nextjs
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/nextjs-stdout.log
stderr_logfile=/var/log/supervisor/nextjs-stderr.log
environment=NODE_ENV="production",PORT="3004",HOSTNAME="0.0.0.0"

; ═══════════════════════════════════════════════════════════════════════════════
; Backend Services (Python FastAPI)
; ═══════════════════════════════════════════════════════════════════════════════

[program:orchestrator]
command=/opt/venv/bin/python -m uvicorn orchestrator:app --host 0.0.0.0 --port 8007
directory=/app/backend/compliance
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/orchestrator-stdout.log
stderr_logfile=/var/log/supervisor/orchestrator-stderr.log

[program:sanctions_service]
command=/opt/venv/bin/python -m uvicorn sanctions_service:app --host 0.0.0.0 --port 8004
directory=/app/backend/compliance
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/sanctions-stdout.log
stderr_logfile=/var/log/supervisor/sanctions-stderr.log

[program:monitoring_service]
command=/opt/venv/bin/python -m uvicorn monitoring_rules:app --host 0.0.0.0 --port 8005
directory=/app/backend/compliance
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/monitoring-stdout.log
stderr_logfile=/var/log/supervisor/monitoring-stderr.log

[program:geo_risk_service]
command=/opt/venv/bin/python -m uvicorn geographic_risk:app --host 0.0.0.0 --port 8006
directory=/app/backend/compliance
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/geo-risk-stdout.log
stderr_logfile=/var/log/supervisor/geo-risk-stderr.log

; ═══════════════════════════════════════════════════════════════════════════════
; Group All Services
; ═══════════════════════════════════════════════════════════════════════════════

[group:amttp]
programs=nginx,nextjs,orchestrator,sanctions_service,monitoring_service,geo_risk_service
priority=999
